# deploy/robbo-obibok.service
[Unit]
Description=Robbo Obibok v2 — Discord Chiptune Radio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=boruta
Group=boruta
WorkingDirectory=/home/boruta/robbo-obibok-v2
EnvironmentFile=/home/boruta/robbo-obibok-v2/.env
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/boruta/robbo-obibok-v2/venv/bin/python3 -u -m src
ExecStartPre=/usr/bin/env sh -c 'rm -f /home/boruta/robbo-obibok-v2/obibok.pid'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
