# deploy/robbo-obibok-v2.service
[Unit]
Description=Robbo Obibok v2 — Discord Chiptune Radio (parallel)
Documentation=https://github.com/wiiii653/robbo-obibok-v2

[Service]
Type=simple
WorkingDirectory=/home/boruta/robbo-obibok-v2
EnvironmentFile=/home/boruta/robbo-obibok-v2/.env
ExecStart=/home/boruta/robbo-obibok-v2/venv/bin/python3 -u -m src
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
