# deploy/robbo-obibok-v2.service
[Unit]
Description=Robbo Obibok v2 — Discord Chiptune Radio
Documentation=https://github.com/wiiii653/robbo-obibok-v2
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=boruta
Group=boruta
WorkingDirectory=/home/boruta/robbo-obibok-v2
EnvironmentFile=/home/boruta/robbo-obibok-v2/.env
Environment=PYTHONPATH=/home/boruta/robbo-obibok-v2
ExecStart=/home/boruta/robbo-obibok-v2/venv/bin/python3 -u -m src
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
